# .config/systemd/user/battery-alert.service
[Unit]
Description=Battery level alert

[Service]
Type=oneshot
ExecStart=%h/bin/battery-alert
LogLevelMax=warning
